# package/multimedia/vdr/udev.d/97-vdr.rules
ENV{ACTION}=="add", KERNEL=="dvb0.demux0", RUN+="/etc/init.d/vdr start &"
